# packages/internal-docs/engineering/getting_started/setup/podman.socket
[Unit]
Description=Podman API Socket
Documentation=man:podman-system-service(1)

[Socket]
ListenStream=%t/podman/podman.sock
SocketMode=0660
TriggerLimitIntervalSec=0
TriggerLimitBurst=0

[Install]
WantedBy=sockets.target
